# .config/systemd/user/easyeffects.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=easyeffects --gapplication-service
Restart=on-failure
